# deploy/10-usb-insert.rules
KERNEL=="sda", RUN+="/bin/bash /home/pi/script/udev-script.sh"
